# webcam-c920/set-c920-zoom.service
[Unit]
Description=Set Logitech C920 webcam zoom level

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set-c920-zoom
